# cookbooks/dartlang/files/default/etc/profile.d/env-dartlang.sh
#!/bin/sh
#
export PATH=/opt/lib/dart-sdk/bin:$PATH
